# .config/profile.d/00-xdg.sh
# ┌───────────────────────── ─ ─   ─
# │ ~/.config/profile
# └───────────────────────── ─ ─   ─

# https://wiki.archlinux.org/title/XDG_Base_Directory

# ····· set lang
export LANGUAGE="fr_FR:fr"

# ····· XDG specification
export XDG_CONFIG_HOME="$HOME/.config"
export XDG_CACHE_HOME="$HOME/.cache"
export XDG_DATA_HOME="$HOME/.local/share"
