# .config/profile.d/20-path.sh
# ····· my bins !
export PATH=$PATH:$HOME/.local/bin
